# systemd/claude-auto-liveness.timer
[Unit]
Description=claude-auto liveness watchdog — каждые 5 минут

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
