# systemd/dept-dispatcher.timer
[Unit]
Description=dept-dispatcher — раз в 2 минуты

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
