# systemd/dept-ledger-snapshot.service
[Unit]
Description=dept-ledger git snapshot (аудит журнала отдела)

[Service]
Type=oneshot
Environment=PATH=/home/rainor/.nvm/versions/node/v20.19.2/bin:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/opt/projects/active/claude-control/bin/dept-ledger rotate
ExecStart=/opt/projects/active/claude-control/bin/dept-ledger snapshot
