# systemd/dept-ledger-snapshot.timer
[Unit]
Description=dept-ledger snapshot — раз в час

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
