# systemd/dept-rebase-check.timer
[Unit]
Description=dept-rebase-check — раз в час

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
